# customize/root/checkkerberos.service
[Unit]
Description=ensure kerberos cam is running

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/startkerberos
